# jcp-gestioninmobiliaria.service
[Unit]
Description=JCP Gestión Inmobiliaria - Real Estate Management System
After=network.target

[Service]
Type=simple
User=jcp
Group=jcp
WorkingDirectory=/opt/jcp-gestioninmobiliaria
ExecStart=/opt/jcp-gestioninmobiliaria/jcp-gestioninmobiliaria
Restart=always
RestartSec=5
EnvironmentFile=/opt/jcp-gestioninmobiliaria/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/jcp-gestioninmobiliaria/pb_data

[Install]
WantedBy=multi-user.target
